# portainer/rootfs/etc/s6-overlay/s6-rc.d/init-portainer/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Private Add-on: Portainer
# Example init script, runs before any other service
# ==============================================================================

bashio::require.unprotected