# apcupsd2mqtt/rootfs/etc/services.d/apcupsd2mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: apcupsd2mqtt
# Runs apcupsd2mqtt script
# ==============================================================================


bashio::log.info "Starting apcupsd2mqtt..."

exec /usr/bin/run.sh
