# apache2-perl-server/rootfs/etc/services.d/httpd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Apache HTTP Server
# Runs httpd
# ==============================================================================
bashio::log.info "Starting Apache HTTP Server"
exec /usr/bin/run1.sh
