# apcupsd2mqtt/rootfs/etc/services.d/apcupsd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start apcupsd service
# ==============================================================================

MAIN_CONFIG=/etc/apcupsd/apcupsd.conf
CONFIG_PATH=/data/options.json
LOCAL_UPS=$(jq --raw-output ".local_ups" $CONFIG_PATH)

if [[ "$LOCAL_UPS" = true ]] ; then
  UPSNAME=$(jq --raw-output ".upsname" $CONFIG_PATH)
  ONBATTERYDELAY=$(jq --raw-output ".onbatterydelay" $CONFIG_PATH)
  BATTERYLEVEL=$(jq --raw-output ".batterylevel" $CONFIG_PATH)
  MINUTES=$(jq --raw-output ".minutes" $CONFIG_PATH)

  if [[ -f $MAIN_CONFIG ]]; then
    sed -i "s/\#UPSNAME$/UPSNAME $UPSNAME/" $MAIN_CONFIG
    sed -i "s/^ONBATTERYDELAY [0-9]\+/ONBATTERYDELAY $ONBATTERYDELAY/" $MAIN_CONFIG
    sed -i "s/^BATTERYLEVEL [0-9]\+/BATTERYLEVEL $BATTERYLEVEL/" $MAIN_CONFIG
    sed -i "s/^MINUTES [0-9]\+/MINUTES $MINUTES/" $MAIN_CONFIG
    sed -i "s/^ANNOY [0-9]\+/ANNOY 0/" $MAIN_CONFIG
    exec /bin/rm -f /etc/apcupsd/powerfail
    exec /sbin/apcupsd -b -f /etc/apcupsd/apcupsd.conf
  else
    echo "[Error] File $MAIN_CONFIG not found!"
  fi
fi
