# ha2gd/rootfs/etc/services.d/ha2gd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: ha2gd
# Runs ha2gd script
# ==============================================================================
bashio::log.info "Starting ha2gd service"
exec /usr/bin/run.sh