# air-raid-alert-ua/rootfs/etc/s6-overlay/s6-rc.d/raid/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Private Add-on: Air Raid Alert Ukraine
# Runs Air Raid Alert Ukraine
# ==============================================================================

bashio::log.info 'Starting Air Raid Alert Ukraine...'

chmod +x /root/raid

exec /root/raid