# air-raid-alert-ua/rootfs/etc/s6-overlay/s6-rc.d/raid/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Private Add-on: Air Raid Alert Ukraine
# Runs Air Raid Alert Ukraine
# ==============================================================================

bashio::log.info 'Starting Radarr...'

chmod +x /opt/bin/raid

exec /opt/bin/raid /share/air-raid/settings.yml