# ua-alerts-monitor/rootfs/etc/s6-overlay/s6-rc.d/alertsmonitor/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Private Add-on: UA Alerts Monitor
# Runs UA Alerts Monitor
# ==============================================================================

bashio::log.info 'Starting UA Alerts Monitor...'

mkdir -p /share/alertsmonitor

chmod +x /usr/bin/alertsmonitor.sh

exec /usr/bin/alertsmonitor.sh