# ua-alerts-monitor/rootfs/etc/s6-overlay/s6-rc.d/alertsmonitor/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Private Add-on: UA Alerts Monitor
# Runs UA Alerts Monitor
# ==============================================================================

bashio::log.info 'Starting UA Alerts Monitor...'

mkdir -p /share/alertsmonitor

curl -s \
  -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
  -H "Content-Type: application/json" \
  -d '{"state": "0000", "attributes": {"message": "Тестове повідомлення", "critical": "false", "friendly_name": "Радар повідомлення", "icon": "mdi:radar"}}' \
  http://supervisor/core/api/states/sensor.radar_status > /dev/null

chmod +x /usr/bin/alertsmonitor.sh

exec /usr/bin/alertsmonitor.sh