# air-raid-alert-ua/rootfs/etc/services.d/raid/run
#!/usr/bin/env bashio
# ==============================================================================
# Home Assistant Private Add-on: Air Raid Alert Ukraine
# Runs Air Raid Alert Ukraine
# ==============================================================================

bashio::log.info 'Starting Air Raid Alert Ukraine...'

exec /usr/bin/raid /share/air-raid/settings.yml

