# mbusreader/rootfs/etc/services.d/mbusreader/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: mbusreader
# Runs mbusreader script
# ==============================================================================
bashio::log.info "Starting mbusreader"
exec /usr/bin/run1.sh
